# openwrt-yun-core/yunbridge.files
# List of files to install for YunBridge OpenWRT integration
/usr/bin/yunbridge
/etc/init.d/yunbridge
/etc/config/yunbridge-ttyath0
