# openwrt-yun-core/uci-defaults/95-yunbridge-silence-kernel-console
#!/bin/sh
# Silence kernel printk messages on serial console to prevent interference
# with YunBridge protocol on ttyATH0.
#
# The Arduino Yun's AR9331 uses ttyATH0 for both:
#   1. Kernel console (at 250000 baud, set in bootargs)
#   2. YunBridge serial protocol (at 115200 baud)
#
# Even though the baud rates differ, kernel printk messages can corrupt
# the COBS-encoded protocol frames. This script silences the kernel
# console via sysctl to prevent frame corruption.
#
# This script runs once on first boot (uci-defaults).

SYSCTL_CONF="/etc/sysctl.d/99-yunbridge-no-console.conf"

# Create sysctl configuration to silence kernel printk
if [ ! -f "$SYSCTL_CONF" ]; then
    mkdir -p /etc/sysctl.d
    cat > "$SYSCTL_CONF" << 'EOF'
# YunBridge: Silence kernel console messages on ttyATH0
# This prevents printk from corrupting serial protocol frames.
# All four values set to 0:
#   - console_loglevel: messages with priority < 0 go to console (none)
#   - default_message_loglevel: default priority for printk (0 = emergency only)
#   - minimum_console_loglevel: minimum allowed console_loglevel
#   - default_console_loglevel: default console_loglevel at boot
kernel.printk = 0 0 0 0
EOF
    echo "[yunbridge] Created $SYSCTL_CONF to silence kernel console"
fi

# Apply immediately if sysctl is available
if command -v sysctl >/dev/null 2>&1; then
    sysctl -p "$SYSCTL_CONF" >/dev/null 2>&1 || true
fi

# Also add to rc.local as backup (in case sysctl.d is not processed early enough)
RC_LOCAL="/etc/rc.local"
MARKER="# YunBridge: Silence kernel console"

if [ -f "$RC_LOCAL" ] && ! grep -q "$MARKER" "$RC_LOCAL"; then
    # Insert before 'exit 0' if present, otherwise append
    if grep -q "^exit 0" "$RC_LOCAL"; then
        sed -i "/^exit 0/i\\
$MARKER\\
echo 0 > /proc/sys/kernel/printk\\
dmesg -n 1\\
" "$RC_LOCAL"
    else
        cat >> "$RC_LOCAL" << EOF

$MARKER
echo 0 > /proc/sys/kernel/printk
dmesg -n 1
EOF
    fi
    echo "[yunbridge] Added kernel console silencing to $RC_LOCAL"
fi

exit 0
