# openwrt-yun-core/uci-defaults/98-yunbridge-default-serial
#!/bin/sh
# Seed YunBridge serial defaults without overwriting user config.
# Runs once on first boot (uci-defaults).

DEFAULT_SERIAL_PORT="/dev/ttyATH0"

current_serial_port="$(uci -q get yunbridge.general.serial_port 2>/dev/null || true)"

if [ -z "$current_serial_port" ]; then
    # Ensure section exists; this is a no-op if already present.
    uci -q set yunbridge.general=general
    uci -q set yunbridge.general.serial_port="$DEFAULT_SERIAL_PORT"
    uci -q commit yunbridge
fi

exit 0
