# openwrt-mcu-core/uci-defaults/97-mcubridge-default-secret
#!/bin/sh
# Seed McuBridge serial shared secret without overwriting user configuration.
# Runs once on first boot (uci-defaults).

PLACEHOLDER="changeme123"

current_secret="$(uci -q get mcubridge.general.serial_shared_secret 2>/dev/null || true)"

if [ -z "$current_secret" ] || [ "$current_secret" = "$PLACEHOLDER" ]; then
    # Ensure section exists; this is a no-op if already present.
    uci -q set mcubridge.general=general

    secret=""

    if command -v openssl >/dev/null 2>&1; then
        secret="$(openssl rand -hex 32 2>/dev/null | tr -d '\n' | tr 'A-F' 'a-f')"
    fi

    if [ -z "$secret" ] && command -v hexdump >/dev/null 2>&1 && [ -r /dev/urandom ]; then
        secret="$(head -c 32 /dev/urandom | hexdump -v -e '/1 "%02x"')"
    fi

    if [ -z "$secret" ] && command -v od >/dev/null 2>&1 && [ -r /dev/urandom ]; then
        secret="$(od -An -N32 -tx1 /dev/urandom | tr -d ' \n')"
    fi

    if [ -z "$secret" ]; then
        # As a last resort, refuse to set a weak/empty secret.
        # The daemon will fail-safe and refuse to start until rotated.
        exit 0
    fi

    uci -q set mcubridge.general.serial_shared_secret="$secret"
    uci -q commit mcubridge
fi

exit 0
