# mcubridge/uci-defaults/96-mcubridge-flash-protection-paths
#!/bin/sh
# Enforce RAM-backed paths by default (flash protection).
# mqtt_spool_dir is ALWAYS kept under /tmp to avoid flash wear.
# allow_non_tmp_paths only affects file_system_root (File component).
# Runs once on first boot / first install (uci-defaults).

DEFAULT_FILE_ROOT="/tmp/yun_files"
DEFAULT_SPOOL_DIR="/tmp/mcubridge/spool"

if ! command -v uci >/dev/null 2>&1; then
    exit 0
fi

allow_non_tmp_paths="$(uci -q get mcubridge.general.allow_non_tmp_paths 2>/dev/null || echo "0")"

current_file_root="$(uci -q get mcubridge.general.file_system_root 2>/dev/null || true)"
current_spool_dir="$(uci -q get mcubridge.general.mqtt_spool_dir 2>/dev/null || true)"

changed=0

needs_tmp_prefix() {
    local candidate="$1"
    [ -n "$candidate" ] || return 0
    [ "$candidate" = "/tmp" ] && return 1
    case "$candidate" in
        /tmp/*) return 1 ;;
    esac
    return 0
}

case "${allow_non_tmp_paths}" in
    1|true|yes|on|enable|enabled)
        : # file_system_root may be non-/tmp when explicitly allowed
        ;;
    *)
        if needs_tmp_prefix "$current_file_root"; then
            uci -q set mcubridge.general=general
            uci -q set mcubridge.general.file_system_root="$DEFAULT_FILE_ROOT"
            changed=1
        fi
        ;;
esac

if needs_tmp_prefix "$current_spool_dir"; then
    uci -q set mcubridge.general=general
    uci -q set mcubridge.general.mqtt_spool_dir="$DEFAULT_SPOOL_DIR"
    changed=1
fi

if [ "$changed" -eq 1 ]; then
    uci -q commit mcubridge
fi

exit 0
